# supervisord.conf
[supervisord]
nodaemon=true

[program:display]
command=/bin/bash /startx.sh
startsecs = 15
